# apiServer.conf
[program:apiServer]
command=/root/anaconda3/bin/gunicorn -c gunicorn.py wsgi:application
directory=/home/backend/
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true